# src/Services/BrowserUseAgent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx24 -nolisten tcp
autorestart=true
priority=100

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=200

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -rfbport 5900 -nopw
autorestart=true
priority=300

[program:websockify]
command=/usr/bin/websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=400

[program:chromium]
# Wipe the profile directory before each (re)start so every launch — including
# supervisord auto-restarts after the sidecar kills Chromium between sessions —
# begins from a clean browser state.
command=/bin/sh -c "rm -rf /tmp/chrome-profile && exec /usr/bin/chromium --no-sandbox --disable-dev-shm-usage --remote-debugging-port=9222 --remote-debugging-address=127.0.0.1 --disable-blink-features=AutomationControlled --start-maximized --user-data-dir=/tmp/chrome-profile"
environment=DISPLAY=":99"
autorestart=true
priority=500
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sidecar]
command=/usr/local/bin/uvicorn server:app --host 0.0.0.0 --port 5041
directory=/app
environment=DISPLAY=":99"
autorestart=true
priority=600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
